# deploy/roles/app/templates/[[app_name]]-api.service.j2
[Unit]
Description={{ app_name }} API (Bun/Hono)
After=network.target

[Service]
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ server_current }}/server
EnvironmentFile={{ app_env_dest }}
ExecStart=/root/.bun/bin/bun ./dist/main.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
